# main.cron
0 0 * * * /urs/bin/python /home/de1h/S_07/main.py